# deploy/playbooks/templates/[[app_name]].service.j2
[Unit]
Description=Fanglao Summer Camp 2026 (TanStack Start + Nitro/Bun)
After=network.target

[Service]
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ server_current }}
EnvironmentFile=/etc/{{ app_name }}.env
ExecStart=/root/.bun/bin/bun .output/server/index.mjs
Environment=NODE_ENV=production
Environment=PORT={{ port }}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
